# rc.local
#!/bin/sh -e
#
# rc.local
#


# Print the IP address
_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  printf "My IP address is %s\n" "$_IP"
fi

#sudo systemctl stop ntp;
#sudo ntpd -q -g >> /root/time.log

sudo rfkill block bluetooth;

sudo timedatectl;

LOC='/root/BBSensor'

cd $LOC && git pull >> /root/git.log; ;

cd $LOC && sudo python3 -m sensorpi >> /root/sensor.log





#echo 'bbsensor00' | sudo tee  /etc/hostname

exit 0
