# rc.local
#!/bin/sh -e
#
# rc.local
#


# Print the IP address
_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  printf "My IP address is %s\n" "$_IP"
fi

sudo systemctl stop ntp
sudo ntpd -q -g


#cd /home/pi/BBSensor && sudo python3 -c 'import opczero;opczero.run()'  >> /root/sensor.log &



exit 0
